# deploy/trading-agent.service
# systemd alternative to docker (see DEPLOY.md option B).
# The agent drives twak via the CLI (reliable approval+swap handling); the
# wallet password comes from TWAK_WALLET_PASSWORD in the .env (EnvironmentFile).
[Unit]
Description=BNB Hack autonomous trading agent (BSC, self-custody via TWAK)
After=network-online.target
Wants=network-online.target

# Replace <user> and the paths below with your deployment user and clone path.
[Service]
Type=simple
User=<user>
WorkingDirectory=/opt/bnb-hack
EnvironmentFile=/opt/bnb-hack/.env
ExecStart=/opt/bnb-hack/.venv/bin/python -m agent --live
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1
KillSignal=SIGTERM
TimeoutStopSec=90

[Install]
WantedBy=multi-user.target
